# launchd/local.wifi-chatwork-notifier.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  これはテンプレートです。 __PROJECT_DIR__ を実パス
  (例: /Users/yourname/workspace/wifi-chatwork-notifier) に置換してから
  ~/Library/LaunchAgents/local.wifi-chatwork-notifier.plist にコピーしてください。
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>local.wifi-chatwork-notifier</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>cd __PROJECT_DIR__ &amp;&amp; npm start</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__PROJECT_DIR__</string>

    <!-- ネットワーク変化のシステム通知に反応して発火する。
         ポーリングも固定スケジュールも使わない、純イベント駆動。 -->
    <key>LaunchEvents</key>
    <dict>
        <key>com.apple.notifyd.matching</key>
        <dict>
            <key>com.apple.system.config.network_change</key>
            <dict>
                <key>Notification</key>
                <string>com.apple.system.config.network_change</string>
            </dict>
        </dict>
    </dict>

    <!-- 起動直後の発火は不要 -->
    <key>RunAtLoad</key>
    <false/>

    <!-- 連続発火を抑制（最低 30 秒は再起動しない） -->
    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>__PROJECT_DIR__/logs/stdout.log</string>

    <key>StandardErrorPath</key>
    <string>__PROJECT_DIR__/logs/stderr.log</string>
</dict>
</plist>
